# config/x-idp.target
[Unit]
Description=X-Net SSO
WantedBy=multi-user.target